# net_services/plugins/tr181-upnp/files/etc/uci-defaults/99-set-uci-config
#!/bin/sh

uci set upnpd.config.enabled='1'
uci set upnpd.config.igdv1='0'

exit 0
