# net_services/plugins/ssh-server/files/etc/uci-defaults/99-disable-default-dropbear
#!/bin/sh

/etc/init.d/dropbear disable
rm -f /etc/config/dropbear

exit 0
