# net_core/plugins/ethernet-manager/files/etc/uci-defaults/99-remove-uci-interfaces
#!/bin/sh

while uci -q delete network.@device[0]; do :; done
uci commit

exit 0