# net_hl_services/plugins/tr181-mqttbroker/files/etc/uci-defaults/99-disable-default-mosquitto-init-script
#!/bin/sh

/etc/init.d/mosquitto disable
rm /etc/init.d/mosquitto

exit 0
