# net_services/plugins/tr181-syslog/files/etc/uci-defaults/99-disable-defaults-syslog-ng
#!/bin/sh

/etc/init.d/syslog-ng disable
rm /etc/init.d/syslog-ng

exit 0
