# prplos/network/utilities/tr181-rest-api/files/etc/uci-defaults/99-amx-fcgi-fixme
#!/bin/sh

# fix amx-fcgi fcgi socket handling
sed -i 's|LIGHTTPD_CHROOT=/webui|LIGHTTPD_CHROOT=/|' /etc/init.d/amx-fcgi

exit 0
