# net_core/plugins/ip-manager/files/etc/uci-defaults/99-remove-uci-ip-configurations
#!/bin/sh

uci delete network.lan 
uci delete network.wan 
uci delete network.wan6
uci commit

exit 0
