# net_services/mods/mod-httpaccess-lighttpd/files/etc/uci-defaults/99-disable-default-lighttpd
#!/bin/sh

/etc/init.d/lighttpd disable

exit 0
