# net_core/plugins/dhcpv4-manager/files/etc/uci-defaults/99-disable-default-dnsmasq-init-script
#!/bin/sh

if [[ -f /etc/init.d/dnsmasq ]]; then
    /etc/init.d/dnsmasq disable
    rm -f /etc/init.d/dnsmasq
fi

exit 0
