# net_core/plugins/dhcpv4-manager/files/etc/uci-defaults/99-remove-dnsmasq-hotplug-script
#!/bin/sh

rm -f /etc/hotplug.d/ntp/25-dnsmasqsec

exit 0
