# net_core/plugins/tr181-routeradvertisement/files/etc/uci-defaults/99-disable-odhcpd-ra-message
#!/bin/sh

uci set dhcp.lan.ra=disable
uci commit

exit 0
