# net_core/plugins/tr181-dns/files/etc/uci-defaults/99-set-failover-controller
#!/bin/sh

sed -i 's|failover-controller = ""|failover-controller = "mod-dns-failover"|' /etc/amx/tr181-dns/tr181-dns_vars.odl

exit 0
