# net_core/plugins/tr181-dhcpv6client/files/etc/uci-defaults/99-disable-netifd-starting-odhcp6c-process
#!/bin/sh

uci set network.wan6.proto=none
uci commit network

exit 0
