# net_services/plugins/tr181-dnssd/files/etc/uci-defaults/99-disable-default-umdns-init-script
#!/bin/sh

/etc/init.d/umdns disable
rm -f /etc/init.d/umdns

exit 0
